# deploy/systemd/eutherbooks-clean.service
[Unit]
Description=Clean old EutherBooks generated audio

[Service]
Type=oneshot
User=nichlas
Group=nichlas
WorkingDirectory=/home/nichlas/EutherBooks
ExecStart=/usr/bin/python3 /home/nichlas/EutherBooks/scripts/clean_generated_audio.py --data-dir /home/nichlas/EutherBooks/data --audio-dir /srv/eutherbooks/audio --library-dir /home/nichlas/EutherBooks/library --keep-per-chapter 1 --floor-keep-per-chapter 1 --failed-keep 0 --min-age-hours 24 --max-audio-gib 12 --backup-keep 3 --backup-max-age-days 14 --clean-python-cache --apply
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
