# deploy/systemd/eutherbooks-clean.service
[Unit]
Description=Clean old EutherBooks generated audio

[Service]
Type=oneshot
User=nichlas
Group=nichlas
WorkingDirectory=/home/nichlas/EutherBooks
ExecStart=/usr/bin/python3 /home/nichlas/EutherBooks/scripts/clean_generated_audio.py --data-dir /home/nichlas/EutherBooks/data --keep-per-chapter 2 --failed-keep 0 --min-age-hours 24 --apply
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
